# contrib/systemd/novaboot-target-off@.service
[Service]
Type=oneshot
ExecStart=/bin/su -c "/home/%i/novaboot-shell off" %i
