# contrib/systemd/novaboot-target-off@.service
[Unit]
Description=Power off a novaboot target %i

[Service]
Type=oneshot
ExecStart=/bin/sh -c "/home/%i/novaboot-shell off"
User=%i
