# server/systemd/novaboot-target-off@.target
[Unit]
Description=Power off a novaboot target %i
Conflicts=novaboot-target-on@%i.service
