# server/systemd/novaboot-server-off@.service
[Unit]
Description=Power off a novaboot target %i

[Service]
Type=oneshot
ExecStart=/usr/bin/novaboot-shell -c off
User=%i
