# server/systemd/novaboot-target-off@.timer
[Timer]
OnActiveSec=10min
Unit=novaboot-target-off@%i.target
