# server/systemd/novaboot-server-off@.service
[Unit]
Description=Power off a novaboot target %i

Conflicts=novaboot-server-on@%i.service

[Service]
Type=oneshot
ExecStart=/bin/true
