# server/systemd/novaboot-server-login@.service
[Unit]
Description=Novaboot target on/off manager for %i
StopWhenUnneeded=true

Conflicts=novaboot-server-off@%i.timer

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/systemctl start novaboot-server-on@%i.service
ExecStop=/bin/systemctl start novaboot-server-off@%i.timer
