# server/systemd/novaboot-power-off@.timer
[Timer]
OnActiveSec=10min
Unit=novaboot-power-off@%i.service
RemainAfterElapse=no
